# deploy/vulntrack.service
# deploy/vulntrack.service
[Unit]
Description=VulnTrack Pro - Vulnerability Management Platform
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=simple
User=vulntrack
Group=vulntrack
ExecStart=/usr/local/bin/vulntrack
EnvironmentFile=/etc/vulntrack/env
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vulntrack

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/vulntrack
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

[Install]
WantedBy=multi-user.target
